# deployments/systemd/oidx-pg-restore-verify.timer
[Unit]
Description=Weekly OpenIDX Postgres restore-verify

[Timer]
# Weekly (Sunday 04:00 UTC), jittered, Persistent so a missed run fires on boot.
OnCalendar=Sun *-*-* 04:00:00
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
