# deployments/apisix-edge/systemd/oidx-oauth.service
[Unit]
Description=OpenIDX oauth service (:8006)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/cmit/openidx
EnvironmentFile=/home/cmit/.config/oidx/common.env
ExecStart=/home/cmit/oidx-runtime/bin/oidx-oauth-service
Restart=always
RestartSec=5
StandardOutput=append:/tmp/oidx-logs/oauth.log
StandardError=append:/tmp/oidx-logs/oauth.log

[Install]
WantedBy=default.target
